# systemd/bloomaquatics-deploy.service
[Unit]
Description=Bloom Aquatics — deploy if origin/main has a green build
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/svei/bloomaquatics
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/home/svei/bloomaquatics/scripts/deploy.sh
